# irix/bin/chanserv.cron
* * * * * /home/chanserv/Bin/checkX
0 * * * * /home/chanserv/Bin/bkup
